# supervisor/etc/supervisord.d/rsz.ini
[program:rsz]  
user=Aario  
command=/usr/local/bin/rsz -p 8888 -d "/var/lib/data/rsz/luexu.com"
autostart=true  
autorestart=true  
stderr_logfile=/var/log/supervisor/rsz.err.log  
stdout_logfile=/var/log/supervisor/rsz.log  
stopsignal=INT  